# 2_packer-build-aws-ami-images/resources/system/apibuilder.service
[Unit]
Description=API Builder Service

[Service]
ExecStart=/home/azureuser/scripts/apibuilder.sh

[Install]
WantedBy=multi-user.target

