# 2_packer-build-aws-ami-images/resources/system/apiportal.service
[Unit]
Description=API Portal Service

[Service]
ExecStart=/home/ubuntu/resources/scripts/apiportal.sh

[Install]
WantedBy=multi-user.target

