# 2_packer-build-aws-ami-images/resources/system/axwayv7_traceability_agent.service
[Unit]
Description=Traceability Agent Service
After=quickstart-online.target

[Service]
User=azureuser
ExecStart=/home/azureuser/scripts/axwayv7_traceabilityagent.sh

[Install]
WantedBy=multi-user.target

