# 2_packer-build-aws-ami-images/resources/system/cassandra.service
[Unit]
Description=Cassandra Service

[Service]
#User=azureuser
#Group=azureuser
ExecStart=/home/ubuntu/resources/scripts/cassandra.sh

[Install]
WantedBy=multi-user.target

